# penctl/image/supervisord-dockerd.conf
[program:dockerd]
command=dockerd
autostart=true
autorestart=true
stdout_logfile=/var/log/dockerd.log
stderr_logfile=/var/log/dockerd.log
